# supervisord.conf
[supervisord]
nodaemon=true

[program:sniproxy]
command=/usr/sbin/sniproxy -f -c /etc/sniproxy.conf
autorestart=true
